# systemd/wallpaper.service
[Unit]
Description=Wallpaper changer

[Service]
Type=oneshot
ExecStart=__PYTHON__ __SCRIPT__ --next