# systemd/wallpaper.timer
[Unit]
Description=Run wallpaper changer periodically

[Timer]
OnBootSec=60min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
